# ch_gentoo/etc/init.d/00-rw-shutdown
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

LIBRWINIT="/etc/conf.d/00-rw"

. "$LIBRWINIT"

provide_lock=${provide_lock-"/dev/.00-rw-init-depend-lock"}

depend() {
	after net
	need net
	_is_allow || return 0
	# prevent exception after set down network
	[ -e "$provide_lock" ] || return 0
	for net in `rc-update | sed -e 's/.*\(net\.[^ ]*\).*\|.*/\1/'`;
	do
		[ -z "$net" ] && continue
		einfo "ro-/: add '$net' to provide list"
		provide "$net"
	done
}

start() {
	_is_allow || return 0
	# set plug to net
	touch "$provide_lock"
	rc-update -u
	# update users /home
	ebegin "home-tar-unpack: set up"
	(
		home_addr=$(_home_get_addr)
		home_root=$(_home_get_root)
		home_driver=$(_home_get_driver)
		home_timeout=$(_home_get_timeout)
		home_timeout=${home_timeout-"1"}
		mode=$(_home_get_mode)
		einfo "home-tar-unpack: use mode='${mode}', driver='${home_driver}', addr='${home_addr}', root='${home_root}', timeout='${home_timeout}'"
		if [ x"$mode" = x"off" -o -z "$home_addr" -o x"$home_driver" = x"none" ];
		then
			einfo "home-tar-unpack: disabled"
			eend 0
			return 0
		fi
		mac=$(_home_get_hwid "$home_driver")
		# print hwids
		echo "$mac" | while read hwid;
		do
			einfo "home-tar-unpack: use hwid=${hwid}"
		done
		if [ -z "$mac" ];
		then
			ewarn "home-tar-unpack: can't get mac from first eth* card"
			eend 1
			return 1
		fi
		# get users' files
		for user in `cd /home; echo *`;
		do
			case "$mode" in
				on)
					if _home_check_user "$user";
					then
						einfo "home-tar-unpack: skip '$user'"
					else
						R=1
						echo "$mac" | while read hwid;
						do
							sudo -u "$user" sh -c ". '$LIBRWINIT'; _home_untar '$home_addr' '${home_root}/${hwid}_${user}' '/home/${user}' '$home_timeout' '$home_driver'"
							[ $? -eq 0 ] && R=0
						done
						if [ $R -ne 0 ];
						then
							ewarn "home-tar-unpack: failed for '$user'"
						else
							einfo "home-tar-unpack: success for '$user'"
						fi
					fi
					;;
				invert)
					if _home_check_user "$user";
					then
						R=1
						echo "$mac" | while read hwid;
						do
							sudo -u "$user" sh -c ". '$LIBRWINIT'; _home_untar '$home_addr' '${home_root}/${hwid}_${user}' '/home/${user}' '$home_timeout' '$home_driver'"
							[ $? -eq 0 ] && R=0
						done
						if [ $R -ne 0 ];
						then
							ewarn "home-tar-unpack: failed for '$user'"
						else
							ewarn "home-tar-unpack: success for '$user'"
						fi
					else
						einfo "home-tar-unpack: skip '$user'"
					fi
					;;
			esac
		done
		# run rsync_
		_home_check_rsyncd
		if [ $? -eq 0 ];
		then
			/etc/init.d/rsyncd start
		fi
	)
	eend $?
	return 0
}

stop() {
	einfo "ro-/: Try normaly shutdown"
	return 0
}

